# arametrics-test-app/config/arametrics-test-app.cron
# crontab: 0 5 * * * -> at 5AM every day
# crontab: 0 12 * * * -> at 12PM every day
# crontab: 0 20 * * * -> at 8PM every day

* * * * * sh /home/aravinth/www/apps/arametrics-test-app/app/config/arametrics-test-app.sh
