# arametrics-test-app/config/arametrics-test-app.cron
# crontab: * * * * * -> every minute

* * * * * sh /home/aravinth/www/apps/arametrics-test-app/app/config/arametrics-test-app.sh
