# arametrics-clockify-google/config/arametrics-clockify-google.cron
# crontab: 0 5 * * * -> at 5AM every day
# crontab: 0 12 * * * -> at 12PM every day
# crontab: 0 20 * * * -> at 8PM every day

* * * * * /home/aravinth/www/apps/arametrics-clockify-google/app/config/arametrics-clockify-google.sh
0 5 * * * /home/aravinth/www/apps/arametrics-clockify-google/app/config/arametrics-clockify-google.sh
0 12 * * * /home/aravinth/www/apps/arametrics-clockify-google/app/config/arametrics-clockify-google.sh
0 20 * * * /home/aravinth/www/apps/arametrics-clockify-google/app/config/arametrics-clockify-google.sh
